# gateway/searxly-gateway.service
[Unit]
Description=Searxly AI gateway
After=network.target

[Service]
Environment=IONET_KEY=PASTE_ROTATED_IONET_KEY_HERE
Environment=ZEROX_KEY=PASTE_0X_API_KEY_HERE
Environment=APP_TOKEN=PASTE_RANDOM_APP_TOKEN_HERE
Environment=PORT=8787
WorkingDirectory=/opt/searxly-gateway
ExecStart=/usr/bin/node /opt/searxly-gateway/gateway.mjs
Restart=always
RestartSec=2
User=www-data
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
